# wayland/systemd/dunst.service.d/session.conf
# dunst dies on output changes sometimes; D-Bus activation would respawn it
# on the next notification, but that notification loses its history —
# proactive restart is better. (Shipped unit already has PartOf=.)
[Service]
Restart=on-failure
