# wayland/systemd/dropbox.service.d/session.conf
# Tie dropbox to the graphical session instead of default.target (the unit's
# stock [Install]): started at login it can't see the display, and restarting
# it later raced KillMode=process orphans (an orphan child made the fresh
# instance exit 0, which Restart=on-failure ignores). The unit is disabled;
# a graphical-session.target.wants symlink (setup.sh add-wants) starts it.
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
