# wayland/systemd/blueman-applet.service.d/session.conf
# The shipped unit has no session binding or restart policy at all.
# After=waybar: don't register the tray icon into a waybar that isn't up yet.
[Unit]
PartOf=graphical-session.target
After=graphical-session.target waybar.service

[Service]
Restart=on-failure
